# deploy/systemd/tls-audit-maintenance.timer
[Unit]
Description=Run TLS Audit maintenance daily

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true

[Install]
WantedBy=timers.target
