# deploy/systemd/tls-audit-stats.timer
[Unit]
Description=Run TLS Audit stats report daily

[Timer]
OnCalendar=*-*-* 03:27:00
Persistent=true

[Install]
WantedBy=timers.target
